# debian_files/DEBIAN/postinst
#!/bin/bash

# required
set -e

# reload systemctl daemon
systemctl daemon-reload

# enable navcoin-angular on boot
systemctl enable navcoin-express
systemctl start navcoin-express
